# crontab
0 * * * * /path/to/python /path/to/speedchecker.py >> /path/to/log/file/speedchecker.log
